# user-demo/disable-wol/50-wol-disable.rules
SUBSYSTEM=="net", ACTION=="add", DRIVERS=="r8169", RUN+="/opt/udevScripts/wol-disable.sh"
